# futuro/plist-gateway-proposto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Versao proposta do plist gateway OpenClaw com ANTHROPIC_API_KEY persistente.
  Diff vs original: 2 linhas novas dentro de <dict> de EnvironmentVariables (key + string com placeholder).
  Substituir __ANTHROPIC_API_KEY_VALUE__ pelo valor real (lendo de Maestro/.env) ANTES de copiar para ~/Library/LaunchAgents/.
  Procedimento manual em Maestro/futuro/B1-resolucao.md.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gateway</string>

    <key>Comment</key>
    <string>OpenClaw Gateway (v2026.4.2)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>1</integer>
    <key>Umask</key>
    <integer>63</integer>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/opt/node/bin/node</string>
      <string>/opt/homebrew/lib/node_modules/openclaw/dist/index.js</string>
      <string>gateway</string>
      <string>--port</string>
      <string>18789</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jesus/.openclaw/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jesus/.openclaw/logs/gateway.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
    <key>HOME</key>
    <string>/Users/jesus</string>
    <key>TMPDIR</key>
    <string>/var/folders/yz/l7yvsq4s2cx6v8sphytflx8h0000gn/T/</string>
    <key>NODE_EXTRA_CA_CERTS</key>
    <string>/etc/ssl/cert.pem</string>
    <key>NODE_USE_SYSTEM_CA</key>
    <string>1</string>
    <key>PATH</key>
    <string>/opt/homebrew/opt/node/bin:/Users/jesus/.local/bin:/Users/jesus/.npm-global/bin:/Users/jesus/bin:/Users/jesus/.volta/bin:/Users/jesus/.asdf/shims:/Users/jesus/.bun/bin:/Users/jesus/Library/Application Support/fnm/aliases/default/bin:/Users/jesus/.fnm/aliases/default/bin:/Users/jesus/Library/pnpm:/Users/jesus/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OPENCLAW_GATEWAY_PORT</key>
    <string>18789</string>
    <key>OPENCLAW_LAUNCHD_LABEL</key>
    <string>ai.openclaw.gateway</string>
    <key>OPENCLAW_SYSTEMD_UNIT</key>
    <string>openclaw-gateway.service</string>
    <key>OPENCLAW_WINDOWS_TASK_NAME</key>
    <string>OpenClaw Gateway</string>
    <key>OPENCLAW_SERVICE_MARKER</key>
    <string>openclaw</string>
    <key>OPENCLAW_SERVICE_KIND</key>
    <string>gateway</string>
    <key>OPENCLAW_SERVICE_VERSION</key>
    <string>2026.4.2</string>
    <key>ANTHROPIC_API_KEY</key>
    <string>__ANTHROPIC_API_KEY_VALUE__</string>
    </dict>
  </dict>
</plist>
