# legado/MONITOR-FONTES/config/launchd/com.stemmia.monitor-fontes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stemmia.monitor-fontes</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/jesus/Desktop/MONITOR-FONTES/scripts/orquestrador.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>13</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jesus/Desktop/MONITOR-FONTES</string>

  <key>StandardOutPath</key>
  <string>/Users/jesus/Desktop/MONITOR-FONTES/logs/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jesus/Desktop/MONITOR-FONTES/logs/launchd.stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
